# roles/system/files/etc/NetworkManager/dispatcher.d/10-no-search
#!/bin/bash
# This file managed by Ansible

if [ "$2" = "up" ] || [ "$2" = "dhcp4-change" -o "$2" = "dhcp6-change" ]; then
  sed -i '/^search /d' /etc/resolv.conf
fi
